# hosts/mbpm3/key_remap/Library/LaunchAgents/sh.sett.keyremap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>sh.sett.keyremap</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sett/Library/LaunchAgents/sh.sett.keyremap.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/sh.sett.keyremap.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/sh.sett.keyremap.stderr.log</string>
</dict>
</plist>
